# Service/shadowsocks.openrc
#!/sbin/openrc-run

# openrc service-script for Shadowsocks
#
# place in /etc/init.d/
# start on boot: "rc-update add Shadowsocks"
supervisor=supervise-daemon
name="Shadowsocks"
description="Shadowsocks Proxy"

command="/root/Shadowsocks/Shadowsocks"
command_args="-config /root/Shadowsocks/config.json"
pidfile="/run/Shadowsocks.pid"

command_background=true

depend() {
    need net
    after net
}